# login/rootimg/etc/profile.d/timezone.sh
export TZ=$(readlink /etc/localtime | cut -d/ -f5-6)
