# controller/rootimg/etc/systemd/system/replicatorwq.service
[Unit]
Description=Synchronize Files accross a cluster

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/replicator
