# controller/rootimg/etc/systemd/system/replicator.timer
[Unit]
Description=Run replicator.service every 1 minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
