# login/rootimg/etc/systemd/system/munge.service.d/customexec.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/munged  --key-file /cluster/etc/munge/munge.key
